# dot/etc/rc.local
#!/bin/sh -e

sudo mount /dev/mmcblk0 /mnt/
mount /dev/mmcblk0 /mnt/
sudo rm -rf ~/.local/share/Trash/*
exit 0
