# dot/etc/rc.local
#!/bin/sh
sudo mount /dev/mmcblk0 /mnt/
