# dot/etc/rc.local
#!/bin/bash -e
mount /dev/mmcblk0 /mnt/sdcard/
exit 0
