# dot/etc/rc.local
#!/bin/sh -e
sudo mount /dev/mmcblk0 /mnt/
exit 0
